# supervisord.conf
[supervisord]
nodaemon=true
loglevel=debug

[program:gunicorn]
command=gunicorn main:app --bind 0.0.0.0:5000
directory=/src/calendar990

[program:nginx]
command=nginx

[program:cron]
command=cron -f
