# supervisord.conf
[supervisord]
nodaemon=true

[program:gunicorn]
command=gunicorn main:app --bind 0.0.0.0:5000
directory=/src/calendar

[program:nginx]
command=nginx
